# ansible/roles/teamcity/templates/teamcity_server.service.j2
#/lib/systemd/system
[Unit]
Description="TeamCity Server"
Requires=network.target
After=syslog.target network.target

[Service]
Type=forking
User={{ tooling_teamcity_svc }}
Group={{ tooling_teamcity_svc }}
EnvironmentFile=/etc/default/teamcity-server.env
PIDFile={{ teamcity_server_dir }}/logs/teamcity.pid
ExecStart={{ teamcity_server_dir }}/bin/teamcity-server.sh start
ExecStop={{ teamcity_server_dir }}/bin/teamcity-server.sh stop

[Install]
WantedBy=multi-user.target
